# supervisord.conf
; =============================================================================
; Foresight EM Agent - Supervisor 进程管理配置
; =============================================================================
; 设计目标：
; 1. 同时管理 Caddy（反向代理）和 Python FastAPI（后端）
; 2. 任一进程崩溃 → 自动重启（autorestart=true）
; 3. 如果重启失败达到上限（FATAL状态）→ event listener 触发容器退出
; 4. Render 检测到容器退出后自动重启整个服务
; =============================================================================

[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/var/run/supervisord.pid
user=root

; ── Python FastAPI 后端 ──────────────────────────────────
[program:backend]
command=uvicorn api.server:app --host 127.0.0.1 --port 8081 --log-level info
directory=/app
user=root
autostart=true
autorestart=true
startretries=5
stopwaitsecs=10
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=
    PYTHONUNBUFFERED="1"

; ── Caddy 反向代理 ──────────────────────────────────────
[program:caddy]
command=caddy run --config /etc/caddy/Caddyfile
user=root
autostart=true
autorestart=true
startretries=5
stopwaitsecs=10
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ── 崩溃检测事件监听器 ──────────────────────────────────
; 当任意进程达到 FATAL 状态（超出重启次数）时，退出 supervisord
; 从而触发容器退出 → Render 自动重启整个服务
[eventlistener:exit_on_crash]
command=python3 /app/exit_on_crash.py
events=PROCESS_STATE_FATAL,PROCESS_STATE_EXITED
buffer_size=1024
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
